# portage/dev-db/apache-cassandra-bin/files/cassandra
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

[ -d /var/run/pids ] || mkdir -m=777 /var/run/pids

pidfile="/var/run/pid/cassandra.pid"
command="/opt/cassandra/bin/cassandra"
command_args="-p $PIDFILE"
start_stop_daemon_args="-u cassandra"


depend() {
    use net
}

